# week2/Virtualisation_and_Linux/crontab
CRONTAB
This is a scheduler it runs set scripts at a certain time. E.g run a script every 30 seconds or every month.

To list all the scheduled cron jobs for the current user, run the following command. It’ll print all the contents of the crontab file.
$ crontab -l

Want the cron job list for a different user? Run this command instead.
$ sudo crontab -u -l

To edit the crontab script, run the command. To edit the crontab file for root, run the following command with sudo privilege.
$ crontab -e

Each line in the crontab script defines a task. Here’s a quick breakdown of the crontab entries.
$ <minute> <hours> <day_of_month> <month> <day_of_week> <command_to_run>
Here’s a list of all the possible values for all these fields. If used an asterisk (*) instead of a numeric value, every possible value of the field will be used.
•	minute: 0 to 59
•	hours: 0 to 23
•	day of the month: 1 to 31
•	month: 1 to 12
•	day of the week: 0 (Sunday) to 6 (Saturday)

Using an asterisk or wildcard in place of any of the five variables will repeat job on all new instances of that variable. The entry below will take a screenshot every minute:
* * * * *  scrot
Every hour:
0 * * * *  scrot
Everyday at 18:00 hours:
0 18 * * * scrot
Every month on 6th at 18:00 hours:
0 18 6 * * scrot
Use “/” to run a script at a fixed periodic interval of 30 minutes (runs every 30 minutes):
*/30 * * * * scrot
You can use comma separated variables to define repeat execution of script at regular intervals. The example below will take a screenshot of your desktop on every 1st, 4th, and 5th minute of a new hour.
1,4,5 * * * *  scrot



One possibility, though it's a bit of a kludge(a), is to have two jobs, one offset by 30 seconds:
# Need these to run on 30-sec boundaries, keep commands in sync.
* * * * *              /path/to/executable param1 param2
* * * * * ( sleep 30 ; /path/to/executable param1 param2 )           (sleep puts the system to sleep for 30 seconds)
You'll see I've added comments and formatted to ensure it's easy to keep them synchronised.
Both cron jobs actually run every minute but the latter one will wait half a minute before executing the "meat" of the job, /path/to/executable.
